# image-files/etc/supervisor/conf.d/export-env.conf
[program:export-env]
command = /usr/local/bin/export-env.sh
autostart = true
autorestart = false
priority=30
startsecs = 0
autorestart = false
startretries = 1
