# image-files/etc/supervisor/conf.d/worker.conf
[program:worker]
command = yii queue/listen
autostart = true
autorestart = true
priority=70
