# image-files/etc/supervisor/conf.d/worker.conf
[program:worker]
command=yii queue/listen
autostart=%(ENV_SUPERVISOR_START_WORKER)s
autorestart=true
priority=70
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
